# com.voice-input.client.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.voice-input.client</string>

    <key>ProgramArguments</key>
    <array>
        <string>/%%PROJECT_ROOT%%/.venv/bin/python</string>
        <string>mac_client.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/%%PROJECT_ROOT%%</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>VOICE_INPUT_SERVER</key>
        <string>ws://localhost:8991</string>
        <key>VOICE_INPUT_HOTKEY</key>
        <string>f13</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/%%USERNAME%%/Library/Logs/voice-input.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/%%USERNAME%%/Library/Logs/voice-input.log</string>

    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
